# luci-theme-saas-dark/modern/root/etc/uci-defaults/99-theme-saas-dark
#!/bin/sh
uci -q batch <<'EOT'
set luci.main.mediaurlbase='/luci-static/saas-dark'
commit luci
EOT
exit 0
