# luci-theme-saas-dark/modern/root/etc/uci-defaults/99-theme-saas-dark
#!/bin/sh

uci -q delete luci.themes.SaaSDark
uci -q set luci.themes.SaaSDark='/luci-static/saas-dark'
uci -q set luci.main.mediaurlbase='/luci-static/saas-dark'
uci -q commit luci

exit 0
