# openwrt-theme-modern/root/etc/uci-defaults/99-theme-modern
#!/bin/sh
uci -q batch <<'EOT'
set luci.main.mediaurlbase='/luci-static/modern'
commit luci
EOT
exit 0
